# dirwatch.service
[Unit]
Description="Simple photo organizer"
Wants=dirwatch.timer

[Service]
Type=oneshot
User=tempdir
Group=photo
ExecStart=/usr/bin/dirwatch.sh

